# examples/minikube/prbot.timer
[Unit]
Description=Minikube performance analysis PR bot
Requires=prbot.service

[Timer]
OnCalendar=*:0/10
